# storm-nimbus/0.9/storm-drpc.conf
[program:storm-drpc]
command=bin/storm drpc
directory=/opt/storm
autorestart=true
