# storm-nimbus/0.9/storm-nimbus.conf
[program:storm-nimbus]
command=bin/storm nimbus
directory=/opt/storm
autorestart=true
