# tomcat/tomcat.conf
[program:tomcat]
command=bin/catalina.sh run
directory=/opt/tomcat
autorestart=true
