# zookeeper/3.4/zookeeper.conf
[program:zookeeper]
command=bin/zkServer.sh start-foreground
directory=/opt/zookeeper
autorestart=true
