# streamflow/0.12/streamflow.conf
[program:streamflow]
command=bin/streamflow.sh
directory=/opt/streamflow
autorestart=true
