# storm-ui/0.9/storm-ui.conf
[program:storm-ui]
command=bin/storm ui
directory=/opt/storm
autorestart=true
