# logstash/logstash.conf
[program:logstash]
command=bin/logstash agent -f /etc/logstash/conf.d
directory=/opt/logstash
autorestart=true
