# storm-supervisor/storm-supervisor.conf
[program:storm-supervisor]
command=bin/storm supervisor
directory=/opt/storm
autorestart=true
